# setup/node-vuln.service
[Unit]
Description=Vulnerability Lookup SSO Application
Documentation=https://github.com/mcleancraig/vuln
After=network.target mariadb.service
Wants=network.target
Requires=mariadb.service

[Service]
Type=simple
User=node-vuln
Group=node-vuln
WorkingDirectory=/opt/node-vuln
ExecStart=/usr/bin/node node-vuln.js
ExecReload=/bin/kill -HUP $MAINPID
KillMode=mixed
KillSignal=SIGTERM
TimeoutStopSec=30
Restart=always
RestartSec=10
StartLimitInterval=60s
StartLimitBurst=3

# Environment
Environment=NODE_ENV=production
Environment=PORT=80
EnvironmentFile=-/etc/node-vuln/environment

# Security hardening
NoNewPrivileges=true
PrivateTmp=true
PrivateDevices=true
ProtectHome=true
ProtectSystem=strict
ReadWritePaths=/opt/node-vuln/logs /tmp
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
RestrictRealtime=true
RestrictNamespaces=true
LockPersonality=true
RestrictAddressFamilies=AF_UNIX AF_INET AF_INET6
SystemCallFilter=@system-service
SystemCallErrorNumber=EPERM

# Needed for port 80 to bind
AmbientCapabilities=CAP_NET_BIND_SERVICE
CapabilityBoundingSet=CAP_NET_BIND_SERVICE

# Resource limits
LimitNOFILE=65536
LimitNPROC=4096
MemoryMax=1G
CPUQuota=200%

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=node-vuln

[Install]
WantedBy=multi-user.target
